# arch/power_profile/update-power-profile.service
[Unit]
Description=Update power profile

[Service]
Type=oneshot
ExecStart=/home/caleb/.local/bin/update-power-profile
